# docker/etc/services.d/tor/run
#!/usr/bin/env bash

set -Eeu
set -o pipefail
export PATH="/usr/sbin:$PATH"


ARGS=(
  -f /srv/run/tor/0-main.conf
  --hush
  )


rm -rf -- /tmp/tor
exec s6-setuidgid "$USER" tor "${ARGS[@]}"

