# docker/etc/services.d/tor/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  -f /srv/run/tor/torrc
  )


#exec s6-setuidgid "$USER" tor "${ARGS[@]}"
sleep infinity

