# docker/etc/services.d/chrony/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  )


exec s6-setuidgid "$USER" chronyd "${ARGS[@]}"
