# docker/etc/services.d/chrony/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  -d
  -f /srv/run/chrony/0-main.conf
  )


chronyd "${ARGS[@]}"
