# docker/etc/services.d/chrony/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  -d
  -x
  -f /srv/run/chrony/0-main.conf
  -L 0
  )


rm -rf -- /tmp/chronyd.pid
exec chronyd "${ARGS[@]}"
