# docker/etc/services.d/nginx/run
#!/usr/bin/env bash

set -Eeu
set -o pipefail
export PATH="/usr/sbin:$PATH"


ARGS=(
  -c /srv/run/nginx/nginx.conf
  )


exec -- s6-envuidgid "$USER" nginx "${ARGS[@]}"
