# docker/etc/services.d/nginx/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  -c /srv/run/nginx/nginx.conf
  )


exec s6-setuidgid "$USER" nginx "${ARGS[@]}"

