# docker/etc/services.d/unbound-cache/run
#!/usr/bin/env bash

set -eu
set -o pipefail


s6-svwait -U /run/s6/legacy-services/unbound

exec s6-setuidgid "$USER" python3 -m router cache
