# docker/etc/services.d/ifup/run
#!/usr/bin/env bash

set -eu
set -o pipefail


while true
do
  python3 -m router ifup
  sleep 1
done
