# docker/etc/services.d/ifup/run
#!/usr/bin/env bash

set -Eeu
set -o pipefail
export PATH="/usr/sbin:$PATH"


while true
do
  s6-envuidgid "$USER" /venv/bin/python3 -m router ifup
  sleep 60
done
