# docker/etc/services.d/redis/run
#!/usr/bin/env bash

set -Eeu
set -o pipefail
export PATH="/usr/sbin:$PATH"


ARGS=(
  -- /srv/run/redis/0-main.conf
  )


exec -- redis-server "${ARGS[@]}"

