# docker/etc/services.d/dnsmasq-dhcp/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  --conf-dir /srv/run/dnsmasq/dhcp/conf.d
  )


s6-svwait -U /run/s6/legacy-services/unbound
exec s6-envuidgid "$USER" dnsmasq "${ARGS[@]}"
