# docker/etc/services.d/unbound/run
#!/usr/bin/env bash

set -eu
set -o pipefail


CONF='/srv/run/unbound/0-include.conf'
ARGS=(
  -c "$CONF"
  )


unbound-checkconf "$CONF"
exec unbound "${ARGS[@]}"
