# docker/etc/services.d/unbound/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  )


exec s6-setuidgid router unbound "${ARGS[@]}"

