# docker/etc/services.d/unbound/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  -d
  -c /srv/run/unbound/0-main.conf
  )


s6-setuidgid "$USER" unbound "${ARGS[@]}"
sleep infinity

