# docker/etc/services.d/dnsmasq-dns/run
#!/usr/bin/env bash

set -Eeu
set -o pipefail
export PATH="/usr/sbin:$PATH"


ARGS=(
  --conf-dir /srv/run/dnsmasq/dns/conf.d
  )


exec -- s6-envuidgid "$USER" dnsmasq "${ARGS[@]}"
