# docker/usr/local/lib/systemd/system/dnsmasq-mon.timer
[Unit]
Description=Restart Dnsmasq DHCP Timer

[Timer]
OnActiveSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=timers.target
