# docker/usr/local/lib/systemd/system/dnsmasq-mon.timer
[Unit]
Description=Restart Dnsmasq DHCP Timer

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
