# docker/usr/local/lib/systemd/system/unbound-mon.timer
[Unit]
Description=Unbound Healthcheck Timer

[Timer]
OnActiveSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=timers.target
