# docker/etc/services.d/unbound-anchor/run
#!/usr/bin/env bash

set -u
set -o pipefail


ARGS=(
  -a /srv/run/unbound/root.key
  )


while true
do
  s6-setuidgid "$USER" unbound-anchor "${ARGS[@]}"
  sleep 3600
done
