# docker/etc/services.d/dnsmasq/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  --conf-dir /srv/run/dnsmasq/conf.d
  )


exec dnsmasq "${ARGS[@]}"
