# docker/etc/services.d/wireguard-webdav/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  serve
  webdav
  --addr=localhost:51820
  --read-only
  /data/wireguard/pub
  )


exec s6-setuidgid "$USER" rclone "${ARGS[@]}"
