# docker/etc/services.d/avahi/finish
#!/usr/bin/env bash

set -Eeu
set -o pipefail


printf -- '%s\n' "FIN :: $PWD"
# exec -- /run/s6/basedir/bin/halt
