# docker/etc/services.d/avahi/finish
#!/usr/bin/env bash

set -eu
set -o pipefail


SVC="$(basename "$0")"
printf '%s' "FIN :: $SVC"
exec s6-svscanctl -t /var/run/s6/services
