# docker/etc/services.d/avahi/finish
#!/usr/bin/env bash

set -eu
set -o pipefail


printf '%s\n' "FIN :: $PWD"
exec /run/s6/basedir/bin/halt
