# docker/etc/services.d/dnsmasq-main/run
#!/usr/bin/env bash

set -eu
set -o pipefail


exec s6-setuidgid router dnsmasq --pid-file=/var/run/dns-main.pid --conf-file=/srv/dnsmasq/0-main.conf
