# docker/etc/services.d/wireguard/run
#!/usr/bin/env bash

set -eu
set -o pipefail


python3 -m router wg
/etc/cont-init.d/00-perms.sh
sleep infinity

