# docker/etc/services.d/wireguard/run
#!/usr/bin/env bash

set -eu
set -o pipefail


#s6-setuidgid "$USER" python3 -m router wg
sleep infinity

