# docker/etc/services.d/cake/run
#!/usr/bin/env bash

set -eu
set -o pipefail


#TODO -- Enable this
#python3 -m router cake
sleep infinity

