# docker/etc/services.d/dnsmasq-lan/run
#!/usr/bin/env bash

set -eu
set -o pipefail


ARGS=(
  --pid-file=/var/run/dns-lan.pid
  --conf-file=/srv/run/dnsmasq/main/0-lan.conf
  )


exec s6-setuidgid router dnsmasq "${ARGS[@]}"
