# docker/etc/services.d/unbound-mon/finish
#!/usr/bin/env bash

set -Eeu
set -o pipefail


exec -- s6-svc -t /run/s6/legacy-services/unbound
