# docker/etc/services.d/unbound-mon/run
#!/usr/bin/env bash

set -eu
set -o pipefail


s6-svwait -U /var/run/s6/services/dnsmasq

while true
do
  s6-setuidgid "$USER" /srv/run/unbound/ctl.sh status &>>/dev/null
  sleep 1
done

