# docker/etc/services.d/unbound-mon/run
#!/usr/bin/env bash

set -eu
set -o pipefail


s6-svwait -U /run/s6/legacy-services/unbound

while true
do
  /etc/services.d/unbound/data/check &>>/dev/null
  sleep 1
done
