# docker/etc/services.d/unbound-mon/run
#!/usr/bin/env bash

set -Eeu
set -o pipefail
export PATH="/usr/sbin:$PATH"


s6-svwait -U /run/s6/legacy-services/unbound

until ! /etc/services.d/unbound/data/check &>>/dev/null
do
  sleep 1
done
